# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/bin/console connect >> /dev/null 2>&1
